# packages/minecraftServerVanilla/etc/systemd/minecraft.service
[Unit]
Description=minecraft server 1.15.1
Documentation=https://minecraft.gamepedia.com/Tutorials/Setting_up_a_server#Linux_instructions
[Service]
WorkingDirectory=/opt/minecraft
Type=simple
ExecStart=/usr/bin/java -Xmx2048M -Xms2048M -jar /opt/minecraft/server.jar nogui
Restart=on-failure
User=minecraft
Group=minecraft
[Install]
WantedBy=multi-user.target